# cygwin64/etc/profile.d/jtcl.sh
#!/bin/bash
export JTCL_HOME=/opt/jtcl/latest
pathadd $JTCL_HOME

export TCLLIBPATH=/opt/testcl/latest

