# cygwin64/etc/profile.d/tools.sh

pathadd /cygdrive/c/tools "custom tools"
alias ls='ls --color'

