# cygwin64/etc/profile.d/zshell.zsh
# set SHELL if not set
if [ ! -z "$ZSH_NAME" -a -z "$SHELL" ]; then
  export SHELL="/bin/zsh"
fi
